# examples/hello/hello.conf
start on socket PROTO=inet PORT=5000

setuid tv
setgid tv
exec /home/tv/go/src/github.com/tv42/upstart-socket/examples/hello/hello
